# deploy/aigc-web.service
[Unit]
Description=AIGC Web Next.js application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/project/AIGC_web
Environment=NODE_ENV=production
Environment=PORT=3010
Environment=HOSTNAME=127.0.0.1
Environment=NEXT_PUBLIC_BASE_PATH=/aigc
ExecStart=/usr/bin/node /home/ubuntu/project/AIGC_web/.next/standalone/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
