# deploy/aigc-storage-cleanup.timer
[Unit]
Description=Run Bala AIGC COS storage cleanup daily

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
