# deploy/aigc-wechat-reconcile.timer
[Unit]
Description=Run WeChat Pay reconciliation every day

[Timer]
OnCalendar=*-*-* 03:20:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
