# deploy/aigc-worker.service
[Unit]
Description=Bala AIGC generation worker
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/project/AIGC_web
Environment=NODE_ENV=production
EnvironmentFile=-/home/ubuntu/project/AIGC_web/.env.production
ExecStart=/usr/bin/node /home/ubuntu/project/AIGC_web/scripts/worker.mjs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
