# deploy/aigc-postgres-backup.timer
[Unit]
Description=Run Bala AIGC PostgreSQL backup daily

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
